# cc.conf
[program:cc]
command=python3 /home/curator2/cc/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/cc/test.err.log
stdout_logfile=/home/curator2/cc/test.out.log